# supervisord.conf
[supervisord]
nodaemon=true
pidfile=/var/run/supervisord.pid
logfile=/var/log/supervisor/supervisord.log

[program:mysql]
command=/etc/init.d/mysql start

[program:http]
command=/start-httpd.sh
